# ocitreeview@.service
[Unit]
Description=OCI Tree View (backend + frontend)
Requires=docker.service
After=docker.service network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=%i
WorkingDirectory=/home/%i/ocitreeview
ExecStart=/usr/bin/docker compose up -d --build
ExecStop=/usr/bin/docker compose down
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
